# systemd/jarvis-pty.service
[Unit]
Description=JARVIS PTY terminal server (:4010)
PartOf=jarvis.target
After=network.target

[Service]
WorkingDirectory=/mnt/e/coding/jarvis-os
EnvironmentFile=-/mnt/e/coding/jarvis-os/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u /mnt/e/coding/jarvis-os/scripts/pty_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=jarvis.target
